# etc/rc.local
zramctl --find --size 4G
mkswap /dev/zram0
swapon /dev/zram0
